# playbooks/app-matlab.socket.j2
[Unit]
Description=gunicorn app-matlab socket

[Socket]
ListenStream=/run/app-matlab.sock

[Install]
WantedBy=sockets.target